# docker/s6-overlay/s6-rc.d/celery-beat/run
#!/command/with-contenv sh
celery -A cata_log.tasks beat --loglevel=${CELERY_LOG_LEVEL:-INFO}
