# docker/s6-overlay/s6-rc.d/migrations/run
#!/command/with-contenv sh
cd /opt/
alembic --config cata_log/migrations/alembic.ini upgrade head
